# cherry-picker/deploy/systemd/cherry-picker-metadata.service
[Unit]
Description=cherry-picker metadata worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cherry-picker
Group=cherry-picker
WorkingDirectory=/opt/cherry-picker
Environment=CHERRY_PICKER_CONFIG=/etc/cherry-picker/metadata.json
ExecStartPre=/opt/cherry-picker/preflight-heat-spool-v2.sh --config /etc/cherry-picker/metadata.json
ExecStart=/opt/cherry-picker/cherry-picker
Restart=always
RestartSec=5
LimitNOFILE=1048576
UMask=0077
StateDirectory=cherry-picker
StateDirectoryMode=0700
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ReadWritePaths=/var/lib/cherry-picker

[Install]
WantedBy=multi-user.target
