# deploy/storage/cherry-storage-wal-upload.service
[Unit]
Description=Upload Cherry PostgreSQL WAL archive off host
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/cherry-backup.env
ExecStart=/usr/local/sbin/cherry-storage-wal-upload
Nice=15
IOSchedulingClass=idle
CPUWeight=10
IOWeight=10
UMask=0077
