# deploy/storage/cherry-storage-backup.timer
[Unit]
Description=Weekly Cherry storage base backup

[Timer]
OnCalendar=Sun *-*-* 03:30:00 UTC
Persistent=true
RandomizedDelaySec=30m
AccuracySec=1m
Unit=cherry-storage-backup.service

[Install]
WantedBy=timers.target
