# go/cherry-picker/deploy/systemd/cherry-picker-discovery.service
[Unit]
Description=cherry-picker discovery worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cherry-picker
Group=cherry-picker
WorkingDirectory=/opt/cherry-picker
Environment=CHERRY_PICKER_CONFIG=/etc/cherry-picker/discovery.json
ExecStart=/opt/cherry-picker/cherry-picker
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target